# scripts/deploy/systemd/docs.service
[Unit]
Description=Echoes Documentation
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=aclarke
WorkingDirectory=/var/www/docs/current
ExecStart=/usr/bin/python3 -m http.server 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
